# config/ufad.service.in
[Unit]
Description=User File Attributes (UFA) Daemon

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/ufad
Type=forking

[Install]
WantedBy=default.target